# docker/rootfs/etc/services.d/oauth2-proxy/run
#!/bin/sh
set -e
set -u

exec oauth2-proxy "$@"
