# docker/rootfs/etc/services.d/oauth2-proxy/run
#!/bin/sh
set -e
set -u

# shellcheck source=params
. ./params

exec oauth2-proxy ${OAUTH2_PROXY_PARAMS} --http-address=unix:/tmp/oauth2-proxy.sock
